# package/openclash-guard/files/etc/init.d/openclash_watchdog

#!/bin/sh /etc/rc.common
START=99

start() {
  echo "$(date '+%Y-%m-%d %H:%M:%S') openclash_watchdog 启动脚本开始执行" >> /tmp/openclash_watchdog_boot.log
  sleep 20
  /bin/sh /root/openclash_watchdog.sh &
}

